# scripts/clawbot/systemd/binance-wg-route.service
# Install to /etc/systemd/system/ via install_binance_wg_route_timer.sh (replaces __BLACKBOX_REPO__).
# Refreshes Binance API IPs in wg-proton-mx when CDN/DNS drifts — see VPN/README.md.

[Unit]
Description=Blackbox — route Binance API (api.binance.com) via wg-proton-mx
Documentation=file:__BLACKBOX_REPO__/VPN/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash __BLACKBOX_REPO__/scripts/clawbot/binance_api_route_via_proton_wg.sh
# Requires CAP_NET_ADMIN — run as root (default for system units).
