# scripts/clawbot/systemd/binance-wg-route.service
# Install to /etc/systemd/system/ via install_binance_wg_route_timer.sh (replaces __BLACKBOX_REPO__).
# Refreshes Binance API IPs in wg-proton-mx when CDN/DNS drifts — see VPN/README.md.

[Unit]
Description=Blackbox — Binance ping knock; WG repair if ping fails
Documentation=file:__BLACKBOX_REPO__/VPN/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Knock (GET /api/v3/ping only); full wg/route repair only if ping != 200 — see binance_api_knock_then_repair_if_needed.sh
ExecStart=/bin/bash __BLACKBOX_REPO__/scripts/clawbot/binance_api_knock_then_repair_if_needed.sh
# Requires CAP_NET_ADMIN for repair path — run as root (default for system units).
