# scripts/clawbot/systemd/binance-wg-route.timer
# Install to /etc/systemd/system/ via install_binance_wg_route_timer.sh
# Every 1 minute: lightweight Binance ping ("knock"); repair only on failure (CDN drift / 451).

[Unit]
Description=Timer — Binance API keepalive knock + WG repair if needed (1 min)

[Timer]
OnBootSec=45s
OnCalendar=*:0/1
Persistent=true
Unit=binance-wg-route.service

[Install]
WantedBy=timers.target
